# ops_package/04_templates/systemd/rotate-logs.timer
[Unit]
Description=Weekly MES log rotation

[Timer]
OnCalendar=Sun *-*-* 03:10:00
Persistent=true

[Install]
WantedBy=timers.target
